# var/lib/dpkg/info/runit.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/runit.service
/usr
/usr/share
/usr/share/runit
/usr/share/runit/inittab
/usr/share/doc
/usr/share/doc/runit
/usr/share/doc/runit/examples
/usr/share/doc/runit/examples/start-stop-daemon.runit
/usr/share/doc/runit/NEWS.Debian.gz
/usr/share/doc/runit/README.Debian
/usr/share/doc/runit/utmpset.8.html
/usr/share/doc/runit/useinit.html
/usr/share/doc/runit/usedietlibc.html
/usr/share/doc/runit/upgrade.html
/usr/share/doc/runit/svlogd.8.html
/usr/share/doc/runit/sv.8.html
/usr/share/doc/runit/runsvdir.8.html
/usr/share/doc/runit/runsvchdir.8.html
/usr/share/doc/runit/runsv.8.html
/usr/share/doc/runit/runscripts.html
/usr/share/doc/runit/runlevels.html
/usr/share/doc/runit/runit.8.html
/usr/share/doc/runit/runit-init.8.html
/usr/share/doc/runit/replaceinit.html
/usr/share/doc/runit/install.html
/usr/share/doc/runit/index.html
/usr/share/doc/runit/faq.html
/usr/share/doc/runit/dependencies.html
/usr/share/doc/runit/chpst.8.html
/usr/share/doc/runit/benefits.html
/usr/share/doc/runit/THANKS
/usr/share/doc/runit/README
/usr/share/doc/runit/changelog.Debian.gz
/usr/share/doc/runit/changelog.gz
/usr/share/doc/runit/copyright
/usr/share/doc/runit/debian
/usr/share/doc/runit/debian/ctrlaltdel
/usr/share/doc/runit/debian/3
/usr/share/doc/runit/debian/2
/usr/share/doc/runit/debian/1
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/utmpset.8.gz
/usr/share/man/man8/update-service.8.gz
/usr/share/man/man8/svlogd.8.gz
/usr/share/man/man8/sv.8.gz
/usr/share/man/man8/runsvdir.8.gz
/usr/share/man/man8/runsvdir-start.8.gz
/usr/share/man/man8/runsvchdir.8.gz
/usr/share/man/man8/runsv.8.gz
/usr/share/man/man8/runit.8.gz
/usr/share/man/man8/runit-init.8.gz
/usr/share/man/man8/chpst.8.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/sv
/usr/sbin
/usr/sbin/update-service
/usr/sbin/runsvdir-start
/usr/sbin/utmpset
/usr/sbin/runsvchdir
/usr/bin
/usr/bin/svlogd
/usr/bin/sv
/usr/bin/runsvdir
/usr/bin/runsv
/usr/bin/chpst
/sbin
/sbin/runit-init
/sbin/runit
/var
/var/lib
/var/lib/supervise
/etc
/etc/sv
/etc/sv/getty-5
/etc/sv/getty-5/finish
/etc/sv/getty-5/run
/etc/service
/etc/sv/getty-5/supervise
